# deb/seed.socket
[Unit]
Description=seed socket

[Socket]
ListenStream=443
NoDelay=true
